# attendance-sync-agent/scripts/attendance-sync-agent.service
[Unit]
Description=BHSPL Attendance Background Synchronization Agent
After=network.target mysql.service

[Service]
Type=simple
WorkingDirectory={INSTALL_DIR}
ExecStart=java -jar {INSTALL_DIR}/target/attendance-sync-agent.jar
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
